# ibisibi.service
[Unit]
Description=ibisibi
Requires=
After=

[Install]
WantedBy=multi-user.target

[Service]
User=pi
Type=simple
ExecStart=/home/pi/bin/ibisibi cycle --from 0 --to 1 --serial /dev/ttyUSB0 --interval-secs 8
Restart=always
RestartSec=10
